# deploy/supervisor/bar.conf
[program:bar]
directory=/
command=/bin/bash -lc 'cd "$BAR_PROJECT_DIR" && exec /bin/bash "$BAR_PROJECT_DIR/deploy/gunicorn_start.sh"'
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/bar/supervisor-out.log
stderr_logfile=/var/log/bar/supervisor-err.log
environment=PYTHONUNBUFFERED="1",BAR_PROJECT_DIR="/srv/bar/app",BAR_VENV_DIR="/srv/bar/app/.venv",BAR_ENV_FILE="/etc/bar/bar.env"
