# deploy/supervisor/bar.conf
[program:__SERVICE_NAME__]
directory=/
command=/bin/bash -lc 'cd "__APP_DIR__" && exec /bin/bash "__APP_DIR__/deploy/gunicorn_start.sh"'
user=__SERVICE_USER__
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=__LOGS_DIR__/out.log
stderr_logfile=__LOGS_DIR__/err.log
environment=LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",PYTHONUNBUFFERED="1",BAR_PROJECT_DIR="__APP_DIR__",BAR_VENV_DIR="__VENV_DIR__",BAR_ENV_FILE="__ENV_FILE__"
